# systemd/vpk-tunnel.service
[Unit]
Description=Cloudflare Quick Tunnel for botvpk
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
ExecStartPre=/bin/bash -c 'truncate -s 0 /tmp/cf_tunnel.log'
ExecStart=/usr/bin/cloudflared tunnel --url http://127.0.0.1:8081 --no-autoupdate
StandardOutput=append:/tmp/cf_tunnel.log
StandardError=append:/tmp/cf_tunnel.log
Restart=always
RestartSec=10
ExecStartPost=-/bin/bash -c 'sleep 5 && /opt/botvpk/scripts/notify-miniapp-url.sh'

[Install]
WantedBy=multi-user.target
